# plex-redirect-stdout/root/etc/services.d/plex-media-scanner-redirect/run
#!/usr/bin/with-contenv bash

umask 022

exec \
        s6-setuidgid root bash /etc/docker-mod-scripts/plex-media-scanner-log-redirect.sh
